# legallens.service.example
[Unit]
Description=Legallens India Next.js Application
After=network.target

[Service]
Type=simple
User=legallens
WorkingDirectory=/home/legallens/legallens-app
Environment="NODE_ENV=production"
Environment="PATH=/home/legallens/.nvm/versions/node/v18.18.0/bin"
ExecStart=/home/legallens/.nvm/versions/node/v18.18.0/bin/npm start
Restart=always
RestartSec=10
StandardOutput=append:/var/log/legallens/app.log
StandardError=append:/var/log/legallens/error.log

[Install]
WantedBy=multi-user.target
